# spindle.service
[Unit]
Description=Spindle daemon for Threadmill
After=default.target

[Service]
Type=simple
ExecStart=%h/.local/bin/spindle
Restart=on-failure
RestartSec=2
Environment=RUST_LOG=info

[Install]
WantedBy=default.target
